# recipes-evo/evo-envinit/evo-envinit/commod-mx6ul/init-gpio.service
[Unit]
Description=Systemd Service to Initialize gpios needed by Communication Module.

[Service]
Type=simple
ExecStart=/bin/bash  /sbin/initgpio.sh

[Install]
WantedBy=multi-user.target
