# recipes-evo/evo-envinit/evo-envinit/mx6ul-comm-module/systemd-firstboot.service
[Unit]
Description=Systemd Service to fix and generate ssh keys.

[Service]
Type=simple
ExecStart=/bin/bash  /sbin/systemd-firstboot.sh

[Install]
WantedBy=multi-user.target
